# files/etc/init.d/gpioled
#!/bin/sh /etc/rc.common

USE_PROCD=1
START=99
STOP=10

start_service() {
	/usr/bin/hgled -r
}

stop_service() {
	/usr/bin/hgled -s
}

restart_service() {
	stop_service
	start_service
}