# files/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

# Auto mount disk drive on startup
# Auto Sinkron Waktu #
#/sbin/sync_time.sh bug.com
exit 0
