# files/etc/rc.local
# Put your custom commands here that should be executed once
# auto sinkron jam dengan bug #
# ketik di terminal curl -i domainmu.com #
# Jika ada respon date tanpa exp berarti bisa di pakai # ganti dengan domainmu yang bawah #
# */15 * * * * /etc/init.d/openclash start #
#/usr/bin/jam.sh www.site.com

exit 0
